# config/supervisor/conf.d/bootup.conf
[program:bootup]
command=/bin/bash -c '/etc/supervisor/bootup.sh'
autostart=true
autorestart=true

