# config/supervisor/conf.d/init-runner.conf
[program:init-runner]
command=/bin/bash -c '/etc/supervisor/run-init-sw.sh'
autostart=true
autorestart=true

